# batbeepd
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opts="start stop restart"

BATTERY_DAEMON=/usr/bin/batbeep
BATTERY_PIDFILE=/var/run/batbeep.pid

start() {
	ebegin "Starting batbeep daemon"
	start-stop-daemon --start --exec "${BATTERY_DAEMON}" --pidfile "${BATTERY_PIDFILE}"
	eend $?
}

stop() {
	ebegin "Stopping batbeep daemon"
	start-stop-daemon --stop --pidfile "${BATTERY_PIDFILE}"
	eend $?
}
